# bytessrvd.conf
#upstart script
start on runlevel [2345]
stop on runlevel [016]


setuid bytessrvd
pre-start exec /opt/bytessrvd/start.sh
pre-stop exec /opt/bytessrvd/stop.sh
